# container/config/supervisord.conf
; ============================================
; Supervisor 配置文件
; Supervisor Configuration File
; ============================================
;
; Supervisor 是一个进程管理器，用于：
; Supervisor is a process manager used to:
;
; 1. 同时启动多个服务 / Start multiple services simultaneously
; 2. 监控服务状态 / Monitor service status
; 3. 服务崩溃时自动重启 / Auto-restart when service crashes
; 4. 管理日志输出 / Manage log output
;
; ============================================

; --------------------------------------------
; 主配置 / Main Configuration
; --------------------------------------------
[supervisord]
; 前台运行（Docker 容器需要）
; Run in foreground (required for Docker containers)
; 如果设置为 false，supervisord 会在后台运行，容器会立即退出
; If set to false, supervisord runs in background, container exits immediately
nodaemon=true

; 日志文件位置
; Log file location
logfile=/var/log/supervisor/supervisord.log

; PID 文件位置
; PID file location
pidfile=/var/run/supervisord.pid

; 运行用户
; Run as user
user=root

; --------------------------------------------
; MariaDB 服务 / MariaDB Service
; --------------------------------------------
[program:mariadb]
; 启动命令 - 直接使用 mysqld
; Start command - use mysqld directly
command=/usr/sbin/mysqld --user=mysql

; 自动启动
; Auto-start when supervisor starts
autostart=true

; 崩溃后自动重启
; Auto-restart when crashes
autorestart=true

; 标准输出日志
; Standard output log
stdout_logfile=/var/log/supervisor/mariadb.log

; 错误输出日志
; Error output log
stderr_logfile=/var/log/supervisor/mariadb-error.log

; 启动优先级（数字越小越先启动）
; Start priority (lower number starts first)
; 数据库需要先于 Apache 启动
; Database needs to start before Apache
priority=10

; 等待几秒确认启动成功
; Wait seconds to confirm startup
startsecs=5

; --------------------------------------------
; Apache 服务 / Apache Service
; --------------------------------------------
[program:apache2]
; 启动命令
; Start command
; -D FOREGROUND 让 Apache 在前台运行
; -D FOREGROUND makes Apache run in foreground
command=/usr/sbin/apache2ctl -D FOREGROUND

; 自动启动
; Auto-start
autostart=true

; 自动重启
; Auto-restart
autorestart=true

; 标准输出日志
; Standard output log
stdout_logfile=/var/log/supervisor/apache2.log

; 错误输出日志
; Error output log
stderr_logfile=/var/log/supervisor/apache2-error.log

; 启动优先级
; Start priority
; 数据库启动后再启动 Apache
; Start Apache after database
priority=20

; 启动延迟，等待数据库先启动
; Startup delay, wait for database to start first
startsecs=3

; --------------------------------------------
; Cron 服务 / Cron Service
; --------------------------------------------
[program:cron]
; 启动命令
; Start command
; -f 让 cron 在前台运行
; -f makes cron run in foreground
command=/usr/sbin/cron -f

; 自动启动
; Auto-start
autostart=true

; 自动重启
; Auto-restart
autorestart=true

; 标准输出日志
; Standard output log
stdout_logfile=/var/log/supervisor/cron.log

; 错误输出日志
; Error output log
stderr_logfile=/var/log/supervisor/cron-error.log

; 启动优先级
; Start priority
priority=30
